# startup/zed.conf
description	"Run Zed local server"
respawn
start on runlevel [23]
setuid user
setgid user
exec /usr/local/bin/zed --local /home/user